# scripts/services/auto_brightness/run
#!/bin/sh

exec 2>&1
sv check llizardGUI > /dev/null || exit 1

exec /usr/bin/sp-als-backlight > /dev/null 2>&1
